# docker/webapp/supervisord.conf
[supervisord]
nodaemon=true

[program:bash]
command=/bin/bash -c /opt/init_db.sh

[program:hhvm]
command=sudo service hhvm start

[program:nginx]
command=sudo service nginx start

